# deploy/systemd/inkstone.service
[Unit]
Description=Inkstone API/Worker Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inkstone
Group=inkstone
WorkingDirectory=/var/lib/inkstone
EnvironmentFile=/etc/inkstone/inkstone.env
ExecStart=/usr/local/bin/inkstone-app --mode all
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inkstone

[Install]
WantedBy=multi-user.target
